# server/com.modalcopy.spellcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.modalcopy.spellcheck</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/brown/Desktop/업무/99_기타/01_모달카피라이터/modalcopy/server/app.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/brown/Desktop/업무/99_기타/01_모달카피라이터/modalcopy/server</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/brown/Desktop/업무/99_기타/01_모달카피라이터/modalcopy/server/server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/brown/Desktop/업무/99_기타/01_모달카피라이터/modalcopy/server/server.error.log</string>
</dict>
</plist>

